# debian/neofs-ir.service
[Unit]
Description=NeoFS InnerRing node
Requires=network.target

[Service]
Type=simple
ExecStart=/usr/bin/epicchain-ir --config /etc/neofs/ir/config.yml
User=epicchain-ir
Group=epicchain-ir
WorkingDirectory=/var/lib/neofs/ir
Restart=always
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
